# deployed/deployed/ansible-deployment/roles/django-app/templates/django-app.service.j2
# Systemd service file for {{ app_name }}
# Generated by Ansible on {{ ansible_date_time.iso8601 }}

[Unit]
Description={{ service_description }}
After=network.target
Wants=network.target

[Service]
Type=notify
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}
Environment=DJANGO_SETTINGS_MODULE={{ django_settings_module }}
Environment=SECRET_KEY={{ django_secret_key }}
Environment=DEBUG={{ django_debug | lower }}
Environment=PYTHONPATH={{ app_dir }}
EnvironmentFile={{ env_file }}
ExecStart={{ venv_path }}/bin/gunicorn {{ django_settings_module.split('.')[0] }}.wsgi:application -c {{ gunicorn_config_file }}
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ app_dir }} {{ log_dir }} {{ static_root }} {{ media_root }}
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target